# deploy/systemd/retrosavemanager-backup-retention.service
[Unit]
Description=RetroSaveManager backup retention cleanup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/retrosavemanager-backup-retention
Environment=RSM_REPO_DIR=/opt/retrosavemanager
Environment=RSM_BACKUP_ROOT=/srv/retrosavemanager/backups
Environment=RSM_KEEP_RECENT=4
Environment=RSM_KEEP_DAYS=7
ExecStart=/bin/bash -lc '"$RSM_REPO_DIR"/deploy/scripts/prune-backups.sh --root "$RSM_BACKUP_ROOT" --keep-recent "$RSM_KEEP_RECENT" --keep-days "$RSM_KEEP_DAYS"'
